# roles/tailscale/templates/tailscale-cert-renewal.service.j2
# Tailscale HTTPS cert renewal service
# Generated by Ansible — do not edit manually

[Unit]
Description=Renew Tailscale HTTPS certificates
After=tailscaled.service
Requires=tailscaled.service

[Service]
Type=oneshot
ExecStart=/usr/bin/tailscale cert --cert-file {{ tailscale_cert_dir }}/{{ tailscale_hostname }}.{{ tailscale_tailnet_name }}.crt --key-file {{ tailscale_cert_dir }}/{{ tailscale_hostname }}.{{ tailscale_tailnet_name }}.key {{ tailscale_hostname }}.{{ tailscale_tailnet_name }}
ExecStartPost=/usr/bin/docker exec caddy caddy reload --config /etc/caddy/Caddyfile
