# roles/tailscale/templates/tailscale-cert-renewal.timer.j2
# Tailscale HTTPS cert renewal timer (weekly)
# Generated by Ansible — do not edit manually

[Unit]
Description=Weekly Tailscale HTTPS certificate renewal

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
